# private_Library/LaunchAgents/com.4jp.context-sync.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.4jp.context-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/domus-context-sync</string>
    </array>

    <!-- Event-driven: fire when registry-v2.json changes on disk -->
    <key>WatchPaths</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/Workspace/meta-organvm/organvm-corpvs-testamentvm/registry-v2.json</string>
    </array>

    <!-- Fallback: also run every 30 minutes in case WatchPaths misses an edit -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Background process settings -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Low priority — context sync is not urgent -->
    <key>Nice</key>
    <integer>15</integer>

    <!-- Low I/O priority -->
    <key>LowPriorityIO</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/domus/context-sync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/domus/context-sync-stderr.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:{{ .chezmoi.homeDir }}/.local/bin</string>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
    </dict>

    <!-- Throttle: don't fire more than once per 30 seconds (debounce rapid edits) -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
