# private_Library/LaunchAgents/com.user.mail_automation.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.mail_automation</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{ .chezmoi.homeDir }}/Workspace/universal-mail--automation/run_automation.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{ .chezmoi.homeDir }}/Workspace/universal-mail--automation</string>

    <!-- Run daily at 9:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Disabled by default -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/domus/mail-automation-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/domus/mail-automation-stderr.log</string>
</dict>
</plist>
